# ansible/roles/restic/files/restic.timer
[Unit]
Description=Timer for Restic backups

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true
Unit=restic.service

[Install]
WantedBy=timers.target
